# targetfs/usr/lib/opkg/info/alsa-state.list
/etc/init.d/alsa-state
/etc/asound.conf
