# targetfs/usr/lib/opkg/info/alsa-state.postrm
#!/bin/sh
if test "x$D" != "x"; then
	OPT="-r $D"
else
	OPT=""
fi
if test "$1" = "remove" -o "$1" = "purge"; then
	if ! test -e "/etc/init.d/alsa-state"; then
		update-rc.d $OPT alsa-state remove
	fi
fi
