# targetfs/usr/lib/opkg/info/dropbear.list
/usr/bin/dbclient
/usr/sbin/dropbearmulti
/etc/init.d/dropbear
/usr/sbin/dropbearkey
/usr/sbin/dropbearconvert
/etc/default/dropbear
/usr/sbin/dropbear
