# targetfs/usr/lib/opkg/info/hostap-daemon.list
/etc/init.d/hostapd
/usr/sbin/hostapd
/usr/sbin/hostapd_cli
/etc/hostapd.conf
