# targetfs/usr/lib/opkg/info/matrix-gui-e.list
/etc/init.d/matrix-gui-e
/usr/bin/matrix_guiE
