# targetfs/usr/lib/opkg/info/matrix-gui-e.postrm
#!/bin/sh
if test "x$D" != "x"; then
	OPT="-r $D"
else
	OPT=""
fi
if test "$1" = "remove" -o "$1" = "purge"; then
	if ! test -e "/etc/init.d/matrix-gui-e"; then
		update-rc.d $OPT matrix-gui-e remove
	fi
fi
