# targetfs/usr/lib/opkg/info/netperf.list
/usr/bin/netperf
/usr/sbin/netserver
/etc/init.d/netperf
