# targetfs/usr/lib/opkg/info/psplash-support.list
/usr/bin/psplash-write
/etc/init.d/psplash
/etc/default/psplash
/etc/default/splashfuncs
/usr/bin/splash-write
