# targetfs/usr/lib/opkg/info/psplash-ti.postinst
#!/bin/sh
#!/bin/sh
update-alternatives --install /usr/bin/psplash psplash /usr/bin/psplash.psplash-ti 20
if test "x$D" != "x"; then
	OPT="-r $D"
else
	OPT="-s"
fi
update-rc.d $OPT psplash start 01 S . stop 20 0 1 6 .
