# targetfs/usr/lib/opkg/info/thttpd.postrm
#!/bin/sh
if test "x$D" != "x"; then
	OPT="-r $D"
else
	OPT=""
fi
if test "$1" = "remove" -o "$1" = "purge"; then
	if ! test -e "/etc/init.d/thttpd"; then
		update-rc.d $OPT thttpd remove
	fi
fi
