# targetfs/usr/lib/opkg/info/thttpd.list
/usr/sbin/htpasswd
/srv/www/cgi-bin/ssi
/usr/sbin/syslogtocern
/usr/sbin/thttpd
/etc/init.d/thttpd
/usr/sbin/makeweb
/srv/www/cgi-bin/phf
/srv/www/cgi-bin/redirect
